# hrv4intervals.timer
[Unit]
Description=Sync Google Fit data to intervals every hour

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target

